# bind/supervisord.conf
[program:named]
command=/work/dist/sbin/named -c /work/dist/etc/named.conf -g